# systemd/llm-fetch-proxy.service
[Unit]
Description=LLM fetch proxy — Chromium 131 URL fetch service for Open WebUI filter
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/kali/nixosconfig/scripts/fetch-proxy.py --port 9879 --host 0.0.0.0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
